# deployment/jetson/forgeos-vision.service
[Unit]
Description=ForgeOS Vision Brain (Jetson → Moonraker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/forgeos-n4pro
Environment=PYTHONPATH=/home/jetson/forgeos-n4pro
# Real-time dynamic loop (0.25s ticks). Add --arm only after shop validation.
ExecStart=/usr/bin/python3 -m forgeos.vision.service --moonraker http://192.168.1.178:7125 --interval 0.25 --config /home/jetson/forgeos-n4pro/configs/vision_rig.yaml
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
